# deploy/pr-dashboard-actuals.timer.template
[Unit]
Description=Daily cloud-cost actuals feed for pr-dashboard

[Timer]
OnCalendar=*-*-* 06:15:00
Persistent=true

[Install]
WantedBy=timers.target
